# systemd/moltbook-heartbeat.timer
[Unit]
Description=ZeptoClaw Moltbook Heartbeat Timer
Requires=moltbook-heartbeat.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
