# systemd/workspace-sync.timer
[Unit]
Description=ZeptoClaw workspace auto-sync timer (every 30 min)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
