# systemd/zeptoclaw-gateway.service
[Unit]
Description=ZeptoClaw Gateway Server
After=network.target

[Service]
Type=simple
ExecStart=%h/.zeptoclaw/scripts/gateway-start.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
